# mail/.config/systemd/user/muchsync.timer
[Unit]
Description=Muchsync Timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=default.target
