# freiburg/.config/systemd/user/rss2email.timer
[Unit]
Description=rss2email timer

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15m

[Install]
WantedBy=default.target

